# config/launchd/com.hq.queue-run.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Mac heavy-job runner. Install:
       ln -sf ~/HQ/config/launchd/com.hq.queue-run.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.hq.queue-run.plist -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hq.queue-run</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>~/HQ/scripts/hq-runner.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/tmp/hq-runner.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/hq-runner.log</string>
</dict>
</plist>
