# docker/crontab
# HQ recurring jobs — executed by supercronic inside the hq container.
# 5-field cron syntax (m h dom mon dow); no user column.
# stdout/stderr are captured by supercronic and streamed to container logs.
#
# This container is the master's dispatch entry point of the queue
# architecture: it detects events (new mail / issue activity), files them as
# jobs under ~/HQ/queue, and runs lightweight single-message triage with a
# small local model. Heavy collect/full-sweep jobs are claimed by workers
# (scripts/hq-runner.sh over ssh for remote workers; docker/crontab.worker's
# master-local-worker.sh for the master's own lower-priority local worker, a
# SEPARATE container — see the hq-worker service in deploy/compose.example.yaml).

*/5  * * * * git -C /home/youruser/HQ pull --rebase --autostash origin main

*/5  * * * * /home/youruser/HQ/bin/hq queue scan
*/5  * * * * /home/youruser/HQ/bin/hq queue work --types triage,thread_update --max 3

# Mirror "Next"-labeled issues to Discord threads (open on Next, archive on
# close). Reuses the hq-github Hermes bot token.
*/10 * * * * python3 /home/youruser/HQ/scripts/discord-issue-sync.py

# GTD daily/cleanup housekeeping (ported from GitHub Actions gtd-daily.yml /
# gtd-cleanup.yml). Container TZ is America/New_York, so these run at local
# 8am/2am — no more UTC/DST drift to reason about.
0 8 * * * /home/youruser/HQ/bin/hq task cron-daily
0 2 * * * /home/youruser/HQ/bin/hq task cron-cleanup
