# docker/crontab
# HQ recurring jobs — run by supercronic inside the hq-cron container.
# 5-field cron syntax (m h dom mon dow); no user column. Paths are relative to
# the bind-mounted repo at /hq. stdout/stderr are captured by supercronic and
# streamed to container logs. Container TZ comes from .env (TZ), so the daily
# entries below run at local time with no UTC/DST drift.
#
# This container is the dispatch entry point of the queue: it keeps the tree
# current, detects events (new mail / issue activity), files them as jobs, and
# runs lightweight single-message triage with a small local model. Heavy
# collect jobs are claimed by the hq-worker service (docker/crontab.worker).

# The one and only git-pull loop.
*/5  * * * * git -C /hq pull --rebase --autostash origin main

# Event detection + lightweight triage.
*/5  * * * * /hq/bin/hq queue scan
*/5  * * * * /hq/bin/hq queue work --types triage,thread_update --max 3

# Mirror "Next"-labeled issues to Discord threads (open on Next, archive on
# close). Optional plugin (`hq install discord`) — a no-op if discord.* config
# is unset.
*/10 * * * * /hq/bin/hq discord sync

# GTD daily/cleanup housekeeping.
0 8 * * * /hq/bin/hq cron daily
0 2 * * * /hq/bin/hq cron cleanup

# NOTE: the nightly Forgejo dump backup is NOT run here. deploy/backup-forgejo.sh
# calls `docker exec forgejo forgejo dump`, which needs the docker socket — and
# this container deliberately does not mount it. Install the backup on the HOST
# crontab instead (see deploy/backup-forgejo.sh header):
#   30 3 * * * /path/to/HQ/deploy/backup-forgejo.sh
