# config/systemd/hq-watch.service
# HQ event watcher + small-model triage (entry-point host).
# Install:
#   mkdir -p ~/.config/systemd/user
#   ln -sf ~/HQ/config/systemd/hq-watch.{service,timer} ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now hq-watch.timer
[Unit]
Description=HQ queue scan + triage

[Service]
Type=oneshot
ExecStart=%h/HQ/scripts/hq-watch.sh
